# questdb/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InfluxDB
# Configures NGINX for use with the Chronograf
# ==============================================================================

# Generate Ingress configuration
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    port "^$(bashio::addon.ingress_port)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf
bashio::log.info "Generate Ingress $(bashio::addon.ip_address):$(bashio::addon.ingress_port)"

if bashio::config.true 'ssl'; then
# Generate influxdb protocol SSL configuration
    bashio::var.json \
        interface "$(bashio::addon.ip_address)" \
        port "9449" \
        certfile "$(bashio::config 'certfile')" \
        keyfile "$(bashio::config 'keyfile')" \
        | tempio \
            -template /etc/nginx/templates/influxlineprotocol.gtpl \
            -out /etc/nginx/servers/influxlineprotocol.stream
    bashio::log.info "Generate Ingress $(bashio::addon.ip_address):$(bashio::config 'ssl_port_influxdb')"
    cp /etc/nginx/servers/influxlineprotocol.stream /config
else
    bashio::log.info "No SSL port for influxdb line protocol"
fi

cp /etc/nginx/servers/ingress.conf /config
cp /etc/nginx/nginx.conf /config
