# questdb/rootfs/etc/s6-overlay/s6-rc.d/questdb/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting QuestDB"
bashio::log.info "Running as $(id -un 2>/dev/null) user"

cd /var/lib/questdb
exec "/entrypoint.sh"
